# systemd/services/CanceledOrders.service
[Unit]
Description=Canceled orders service

[Service]
Type=oneshot
ExecStart=
