# systemd/services/CheatedPhones.service
[Unit]
Description=Cheated phones service

[Service]
Type=oneshot
ExecStart=