# systemd/services/NewStopSalesByIngredients.service
[Unit]
Description=New stop sales by ingredients service

[Service]
Type=oneshot
ExecStart=
