# systemd/services/ResetCanceledOrdersStorage.service
[Unit]
Description=Reset canceled orders storage service

[Service]
Type=oneshot
ExecStart=
