# systemd/services/ResetStopSalesByIngredientsStorage.service
[Unit]
Description=Reset stop sales by ingredients storage service

[Service]
Type=oneshot
ExecStart=
