# systemd/services/ResetCheatedPhones.service
[Unit]
Description=Reset cheated phones storage service

[Service]
Type=oneshot
ExecStart=