# systemd/services/StopSalesByPartialIngredients.service
[Unit]
Description=Stop sales by partial ingredients service

[Service]
Type=oneshot
ExecStart=