# systemd/services/StopSalesByAllIngredients.service
[Unit]
Description=Stop sales by all ingredients service

[Service]
Type=oneshot
ExecStart=
