# systemd/services/StopSalesBySectors.service
[Unit]
Description=Stop sales by sectors service

[Service]
Type=oneshot
ExecStart=