# systemd/services/StopSalesBySalesChannels.service
[Unit]
Description=Stop sales by sales channels service

[Service]
Type=oneshot
ExecStart=
