# systemd/services/StopSalesByStreets.service
[Unit]
Description=Stop sales by streets service

[Service]
Type=oneshot
ExecStart=