# systemd/timers/ResetCanceledOrdersStorage.timer
[Unit]
Description=Reset canceled orders storage timer

[Timer]
OnCalendar=00:00:00

[Install]
WantedBy=timers.target
