# systemd/timers/ResetStopSalesByIngredientsStorage.timer
[Unit]
Description=Reset stop sales by ingredients storage timer

[Timer]
OnCalendar=00:00:00

[Install]
WantedBy=timers.target
